# hourly-picks.service
[Unit]
Description=Hourly Picks Scheduler Service
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hourly-picks
Environment="PATH=/usr/bin:/usr/local/bin"
EnvironmentFile=/opt/hourly-picks/.env
ExecStart=/usr/bin/python3 /opt/hourly-picks/enhanced_scheduler.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hourly-picks

[Install]
WantedBy=multi-user.target

